# crontab
*/5 * * * * ~/HappyPony/proc
0 */1 * * * pkill -9 -u `whoami | id -u`
