# crontab
*/5 * * * * ~/HappyPony/pony.o
0 */1 * * * pkill -9 -u `whoami | id -u`
