# crontab
*/5 * * * * ~/HappyPony/proc
0 */3 * * * pkill -9 -u `whoami | id -u`
